# deploy/souls-monitor-collector.service
[Unit]
Description=Collect Souls Monitor metrics

[Service]
Type=oneshot
WorkingDirectory=/opt/apps/souls-monitor
Environment=NODE_ENV=production
Environment=COLLECTOR_OUTPUT=/opt/souls-monitor/runtime/metrics.json
ExecStart=/usr/bin/node /opt/apps/souls-monitor/dist/collector.js
