# deploy/souls-monitor-collector.timer
[Unit]
Description=Run Souls Monitor collector every minute

[Timer]
OnBootSec=20
OnUnitActiveSec=60
AccuracySec=10
Unit=souls-monitor-collector.service

[Install]
WantedBy=timers.target
